# container-files/etc/supervisor.d/services.conf
[program:premount]
command = /mount.sh
autostart = true
autorestart = false
priority = 998

[program:mount]
command = /usr/bin/mount -a
autorestart = true
priority = 999