# container-files/etc/supervisor.d/services.conf
[program:glusterd]
command = /sbin/glusterd
autorestart = true
priority = 998

[program:createvol]
command = /install/start.sh
autostart = true
autorestart = false
priority = 999
